# deployment/systemd/gex-ingestion.service
[Unit]
Description=GEX Options Data Ingestion
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/gex-options-platform/src/ingestion
Environment="PATH=/home/ubuntu/gex-options-platform/venv/bin"
Environment="PYTHONPATH=/home/ubuntu/gex-options-platform/src/ingestion"
EnvironmentFile=/home/ubuntu/gex-options-platform/.env

ExecStart=/home/ubuntu/gex-options-platform/venv/bin/python tradestation_streaming_ingestion_engine.py

# Restart behavior
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gex-ingestion

[Install]
WantedBy=multi-user.target
