# deployment/systemd/gex-ingestion.service
[Unit]
Description=ZeroGEX Options Data Ingestion Engine Service
Documentation=https://github.com/zero-gex-options/gex-options-platform
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/gex-options-platform

Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/bin:/usr/local/bin"

ExecStart=/home/ubuntu/gex-options-platform/deployment/run/ingestion.sh

# Restart behavior
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Resource limits
MemoryMax=2G
MemoryHigh=1.6G
CPUQuota=75%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/gex-options-platform /data /var/log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gex-ingestion

[Install]
WantedBy=multi-user.target
